# static-modules/gogs/templates/gogs.service.erb
[Unit]
Description=Gogs
After=syslog.target
After=network.target

[Service]
# Modify these two values and uncomment them if you have
# repos with lots of files and get an HTTP error 500 because
# of that
###
#LimitMEMLOCK=infinity
#LimitNOFILE=65535
Type=simple
User=<%= @user %>
Group=<%= @user %>
WorkingDirectory=/usr/share/gogs/<%=@version%>/gogs/
ExecStart=/usr/share/gogs/<%=@version%>/gogs/gogs web
Restart=always
Environment=USER=<%= @user %> HOME=/home/<%= @user %>

[Install]
WantedBy=multi-user.target
